# paige-ui.service
[Unit]
Description=PAIGE UI — voice assistant display
After=multi-user.target sound.target
Wants=sound.target

[Service]
Type=simple
User=rpi
WorkingDirectory=/home/rpi/Desktop/PAIGE

# Give the framebuffer and GPIO time to initialise on boot
ExecStartPre=/bin/sleep 3

ExecStart=/home/rpi/Desktop/PAIGE/venv/bin/python /home/rpi/Desktop/PAIGE/app.py

# Kill every process in the cgroup on stop (ensures GPIO pins are freed)
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=5

Restart=on-failure
RestartSec=5

EnvironmentFile=-/home/rpi/Desktop/PAIGE/.env
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
